# musicbot.service
[Unit]
Description=Music Downloader Telegram Bot
After=network-online.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
# Adjust User and paths to match your server setup.
User=ubuntu
WorkingDirectory=/home/ubuntu/Music_Downloader_bot
ExecStart=/home/ubuntu/Music_Downloader_bot/env/bin/python /home/ubuntu/Music_Downloader_bot/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
